# terminal-server.service
[Unit]
Description=Mobile Terminal Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/selstad/Desktop/terminal-workspace/terminal
ExecStart=/bin/bash /home/selstad/Desktop/terminal-workspace/terminal/dev.sh
Restart=always
RestartSec=5
Environment=NODE_ENV=development

[Install]
WantedBy=multi-user.target
